# Linux/Ubuntu/systemd/unit_example/timer.timer
[Unit]
Description=Example

[Timer]
Unit=monitoring.service
# Каждую минуту
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
